# relay/deploy/cyphes-network.service
[Unit]
Description=CYPHES relay and rendezvous network
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cyphes-network
Group=cyphes-network
Environment=CYPHES_RELAY_DATA_DIR=/var/lib/cyphes-network
Environment=CYPHES_RELAY_PORT=4001
EnvironmentFile=-/etc/cyphes-network.env
ExecStart=/usr/local/bin/cyphes-relay
Restart=always
RestartSec=5
StateDirectory=cyphes-network
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
ReadWritePaths=/var/lib/cyphes-network

[Install]
WantedBy=multi-user.target
